# crontab
*/1 * * * * /usr/local/bin/python /root/FogosPtAlerts/FogosPtAlerts.py
